# backend/scripts/velocity-sync.service
[Unit]
Description=Velocity Sync Daemon
After=network.target pocketbase.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/hololive.com.cn/backend/scripts
ExecStart=/usr/bin/node sync_velocity.js
Restart=always
RestartSec=10
Environment=PB_URL=http://127.0.0.1:8090
Environment=PB_EMAIL=admin@local.dev
Environment=PB_PASS=password123456
Environment=VELOCITY_DIR=/opt/velocity

[Install]
WantedBy=multi-user.target
